# .config/systemd/user/krice.service
[Unit]
Description=Krice Plasma Auto Wal

[Service]
ExecStart=%h/.config/krice/krice.sh
Restart=always

[Install]
WantedBy=default.target